# servics.service

[Unit]
Description=bot

[Service]
ExecStart=
WorkingDirectory=
Restart=on-failure

[Install]
WantedBy=multi-user.target
