# xyz.service
[Unit]
Description=bot

[Service]
WorkingDirectory=/root/
ExecStart=/root/
Restart=on-failure

[Install]
WantedBy=multi-user.target
